# infra/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
user=postgres

[program:redis]
command=/usr/bin/redis-server --save 60 1 --appendonly yes
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:api]
directory=/app/apps/api
command=/usr/local/bin/node dist/src/main.js
environment=NODE_ENV="production",PORT="3001",DATABASE_URL="postgresql://invenzo:invenzo@127.0.0.1:5432/invenzo",REDIS_URL="redis://127.0.0.1:6379",JWT_ACCESS_SECRET="invenzo-access-secret",JWT_REFRESH_SECRET="invenzo-refresh-secret"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:web]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
